# crontab
30 2 * * * /bin/bash /app/docker-run.sh
